# automatisation/crontab
52 9 * * 3 /app/start_spider.sh >> /var/log/cron.log 2>&1
